# supervisord.conf
[supervisord]
nodaemon=true

[program:flussonic]
command=/etc/init.d/flussonic start
autorestart=true

